# crontab
0	1	*	*	*	rm -f /tmp/logs/frame.txt /tmp/logs/perform.txt /tmp/logs/homiq.txt /tmp/logs/store.txt /tmp/logs/cal.txt 
#5	0,8,10,16,22	*	*	*	/homekameleon/fsync.sh
