# crontab
0 7 * * * wget -O /dev/null http://gold.webkameleon.com/startmebaby >/dev/null 2>/dev/null
