# promienko
#!/bin/sh /etc/rc.common
# Example script
# Copyright (C) 2007 OpenWrt.org
 
START=70
STOP=75
 
start() {        
	cd /sd/promienko
	export TZ=CET-1CEST,M3.5.0,M10.5.0/3
	./promienko.sh &
}                 
 
stop() {          
	killall -TERM promienko.sh
	kill -TERM `cat /sd/promienko/app.pid` 
	#killall -TERM ssh
}
